# crons/schedule.crontab
# Prescient Data Aggregation Server
0 */12 * * * bash /home/dm1447/dev/ampscz-interview-qc/crons/prescient-transcript-tracker.sh
